# templates/default/systemd.service.erb
[Unit]
Description=Advanced key-value store
After=network.target
Description=<%= @name %>

[Service]
Type=forking
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
ExecStop=/bin/kill -s <%= @stop_signal %> $MAINPID
PIDFile=/var/run/redis/redis-server.pid
TimeoutStopSec=0
Restart=always

User=<%= @user %>
Group=<%= @group %>
WorkingDirectory=<%= @directory %>
Restart=on-failure

PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ReadOnlyDirectories=/
ReadWriteDirectories=-/var/lib/redis
ReadWriteDirectories=-/var/log/redis
ReadWriteDirectories=-/var/run/redis
CapabilityBoundingSet=~CAP_SYS_PTRACE

# redis-server writes its own config file when in cluster mode so we allow
# writing there (NB. ProtectSystem=true over ProtectSystem=full)
ProtectSystem=true
ReadWritePaths=-/etc/redis

[Install]
WantedBy=multi-user.target
Alias=redis.service
